# systemd/dot-config/systemd/user/weather.service
[Unit]
Description=Fetch current weather data.

[Service]
Type=oneshot
ExecStart=%h/.scripts/weather
