# systemd/dot-config/systemd/user/weather.timer
[Unit]
Description=Fetch current weather data.

[Timer]
OnBootSec=1m
OnUnitActiveSec=60m

[Install]
WantedBy=timers.target
