# systemd/.config/systemd/user/mpd.service
[Unit]
Description=Music Player Daemon
Wants=mpdas.service

[Service]
ExecStart=/usr/bin/mpd --no-daemon

[Install]
Also=mpd.socket
