# systemd/.config/systemd/user/shiba.service
[Unit]
Description=Fetch current SHIB price.

[Service]
Type=oneshot
ExecStart=%h/.scripts/shib
