# systemd/.config/systemd/user/shiba.timer
[Unit]
Description=Fetch current SHIB price.

[Timer]
OnActiveSec=60m
OnUnitActiveSec=60m

[Install]
WantedBy=timers.target
