# systemd/.config/systemd/user/sitalarm.timer
[Unit]
Description=Track sitting time

[Timer]
OnActiveSec=45m
OnUnitActiveSec=45m

[Install]
WantedBy=timers.target
