# systemd/.config/systemd/user/wallpaper.timer
[Unit]
Description=Change wallpaper hourly.

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
