# systemd/.config/systemd/user/wallpaper.service
[Unit]
Description=Change wallpaper

[Service]
ExecStart=%h/.scripts/wallchanger.sh
