# cron/etc/cron.weekly/fstrim
#!/bin/sh

# trim all mounted file systems which support it
/usr/bin/fstrim --all || true
